# fmk/rootfs/etc/uci-defaults/luci-theme-40-r3000_base
#!/bin/sh
uci batch <<-EOF
	set luci.themes.r3000_base=/luci-static/r3000_base
        commit luci
EOF
exit 0
	
