# fmk/rootfs/usr/lib/opkg/info/firewall.list
/lib/firewall/core_interface.sh
/lib/firewall/core_ddos.sh
/etc/init.d/firewall
/lib/firewall/config.sh
/etc/config/firewall
/lib/firewall/core.sh
/lib/firewall/fw.sh
/lib/upgrade/keep.d/firewall
/lib/firewall/core_rule.sh
/lib/firewall/core_redirect.sh
/lib/firewall/core_algs.sh
/lib/firewall/core_init.sh
/sbin/fw
/lib/firewall/uci_firewall.sh
