# fmk/rootfs/usr/lib/opkg/info/igmpproxy.list
/etc/init.d/igmp
/etc/config/igmp
/usr/sbin/igmpproxy
/etc/igmpproxy.conf
