# fmk/rootfs/usr/lib/opkg/info/luci-app-userfw.list
/usr/lib/userfw.sh
/usr/lib/lua/luci/model/cbi/userfw/userfw.lua
/etc/init.d/userfw
/etc/config/userfw
/etc/uci-defaults/luci-userfw
/usr/lib/lua/luci/i18n/userfw.zh-cn.lmo
/usr/lib/lua/luci/controller/userfw.lua
