# fmk/rootfs/usr/lib/opkg/info/luci-mod-admin-core.list
/usr/lib/lua/luci/view/admin_status/lan.htm
/etc/init.d/luci_fixtime
/usr/lib/lua/luci/view/admin_status/wan.htm
/usr/lib/lua/luci/view/header.htm
/usr/lib/lua/luci/model/cbi/admin_status/lan.lua
/usr/lib/lua/luci/view/footer.htm
/usr/lib/lua/luci/view/error404.htm
/usr/lib/lua/luci/view/admin_status/help.htm
/usr/lib/lua/luci/model/cbi/admin_status/wan.lua
/usr/lib/lua/luci/view/admin_status/index.htm
/usr/lib/lua/luci/view/acserror.htm
/usr/lib/lua/luci/view/header_help.htm
/usr/lib/lua/luci/tools/webadmin.lua
/usr/lib/lua/luci/controller/admin/servicectl.lua
/usr/lib/lua/luci/tools/status.lua
/usr/lib/lua/luci/view/acslogin_process.htm
/usr/lib/lua/luci/view/sysauth.htm
/usr/lib/lua/luci/model/cbi/admin_services/normal.lua
/etc/init.d/luci_dhcp_migrate
/usr/lib/lua/luci/view/acslogin.htm
/usr/lib/lua/luci/view/error500.htm
/usr/lib/lua/luci/view/styles.htm
/usr/lib/lua/luci/view/admin_status/remote.htm
/usr/lib/lua/luci/view/indexer.htm
