# fmk/rootfs/usr/lib/opkg/info/qostool.list
/etc/init.d/qos
/etc/config/qos
/sbin/qostool
