# fmk/rootfs/usr/lib/opkg/info/rate_limit.list
/etc/config/rateLimit
/etc/init.d/rateLimit
