# fmk/rootfs/usr/lib/opkg/info/tr069.list
/sbin/restore.sh
/etc/config/maintenance
/etc/tr069/tr.xml
/etc/config/remote
/sbin/shpasswd
/etc/tr069/tr.conf
/sbin/tr069
/etc/config/monitor
/sbin/backup.sh
/etc/config/devinfo
/etc/config/userinfo
/etc/config/diagnostics
/etc/tr069/ca.pem
/etc/config/alarm
/sbin/tr_send
/etc/init.d/tr069
/etc/tr069/client.crt
/etc/tr069/client.key
/sbin/devinfo.sh
